# profile.d/ztokens.sh
# 用于存放一些token，如Github API Token，如：
export GITHUB_API_TOKEN=xxx
